# gridpythia.service
[Unit]
Description=GridPythia – Home energy prediction and optimization
Documentation=https://github.com/cnadler86/GridPythia
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gridpythia
Group=pythia
WorkingDirectory={{WORKING_DIR}}
ExecStart={{PYTHON_BIN}} main.py {{EXTRA_ARGS}}

Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging via systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gridpythia

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{WORKING_DIR}} /var/cache/gridpythia-uv
ProtectHome=true
Environment=UV_CACHE_DIR=/var/cache/gridpythia-uv

[Install]
WantedBy=multi-user.target
